# app-portage/portage-systemd-timers/files/4/emerge-user-timers.target
[Unit]
Description=Target for emerge-related user timers.
