# app-portage/portage-systemd-timers/files/4/emerge-user-timer-delay.timer
[Unit]
Description=Enforce a startup delay for emerge-related user timers.

[Timer]
Persistent=False
OnActiveSec=300
Unit=emerge-user-timers.target

[Install]
WantedBy=timers.target
