# app-portage/portage-systemd-timers/files/4/trim-emerge-log.timer
[Unit]
Description=Trim emerge logs to limit overall size

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
