# app-portage/portage-systemd-timers/files/1/eclean-packages.timer
[Unit]
Description=Daily cleanup of old binary packages
RequireMountsFor=/var/cache/binpkgs

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
