# app-emulation/qemu-init-scripts/files/0/qemu
#!/sbin/runscript
# Copyright 1999-2014 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
    use net
}

INSTANCE=${SVCNAME#*.}
conf=$(add_suffix /etc/conf.d/qemu)
pidfile=/run/qemu-${INSTANCE}.pid
monsock=/run/qemu-${INSTANCE}.monitor
source ${conf}

start() {
    ebegin "Starting QEMU ${INSTANCE}"
    start-stop-daemon --pidfile ${pidfile} --background --user qemu --make-pidfile --wait 500 --nicelevel 5 --exec /usr/bin/qemu-system-${arch} -- -name ${INSTANCE},process=qemu-${INSTANCE} -machine ${machine} -cpu ${cpu} -m ${mem} --sandbox on -usb -usbdevice tablet -display vnc -vga std -vnc ${vnc} ${extra_opts} -monitor socket,path=${monsock},server,nowait
    eend $?
}

stop() {
    ebegin "Initiating shutdown of QEMU ${INSTANCE}"
    echo "system_poweroff" > ${monsock}
    sleep ${kill_delay}
    if start-stop-daemon --stop --test --pidfile ${pidfile} --user qemu --exec qemu-system-${arch} ; then
        eend 0
    else
        start-stop-daemon --stop --pidfile ${pidfile} --user qemu --exec qemu-system-${arch}
        eend $?
    fi
}
