# app-portage/portage-systemd-timers/files/3/eclean-distfiles.timer
[Unit]
Description=Daily cleanup of old distfiles
RequiresMountsFor=/var/cache/distfiles

[Timer]
Persistent=true
OnCalendar=daily

[Install]
WantedBy=timers.target
