# app-portage/portage-systemd-timers/files/3/emerge-sync-check.timer
[Unit]
Description=Track whether the system emerge-sync-done target is active or not

[Timer]
Persistent=true
OnCalendar=*:00/15
AccuracySec=900

[Install]
WantedBy=timers.target
