# app-portage/portage-systemd-timers/files/0/eclean-distfiles.service
[Unit]
Description=Clean up old distfiles
RequireMountsFor=/var/cache/distfiles

[Service]
Type=oneshot
User=root
Group=root
ExecStart=eclean -dv distfiles
