# app-portage/portage-systemd-timers/files/0/emerge-sync-done.target
[Unit]
Description=Portage repository synchronization complete
AllowIsolate=No
