# app-portage/portage-systemd-timers/files/0/glsa-notify.target
[Unit]
Description=Send notifications about the system being affected by one or more GLSAs
AllowIsolate=No
